# deployment/mama-time.service.example
[Unit]
Description=MAMA TIME React + Node + PostgreSQL application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mamatime
Group=mamatime
WorkingDirectory=/var/www/mama-time
EnvironmentFile=/var/www/mama-time/.env
ExecStartPre=/usr/bin/npm run db:migrate
ExecStart=/usr/bin/node /var/www/mama-time/backend/src/server.js
Restart=always
RestartSec=5
TimeoutStopSec=15
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/www/mama-time/backend/backups

[Install]
WantedBy=multi-user.target
